# profile.d/github.sh
export GITHUB_TOKEN=9e182793db9c370f7df890f3bf6ad82c71192615

